# launchd/com.slee2.muxtunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.slee2.muxtunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>cd /Users/slee2/projects/muxtunnel &amp;&amp; /etc/profiles/per-user/slee2/bin/npx concurrently "bash -c 'export NAME=mux-service &amp;&amp; exec /etc/profiles/per-user/slee2/bin/node node_modules/.bin/tsx watch src/server.ts'" "bash -c 'export NAME=mux &amp;&amp; exec /etc/profiles/per-user/slee2/bin/node node_modules/.bin/vite'"</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/slee2/Library/Logs/muxtunnel.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/slee2/Library/Logs/muxtunnel.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/slee2</string>
        <key>PATH</key>
        <string>/etc/profiles/per-user/slee2/bin:/run/current-system/sw/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>NAME</key>
        <string>mux</string>
    </dict>
</dict>
</plist>
